# cron/crontab
#crontab

#*/5	*   *   *   *	root 	*applications/medibox/cron/notifications.py
*   *   *   *   *	root	*applications/medibox/cron/test.py
