# cron/crontab
#crontab

*	*   *   *   *	root 	*applications/medibox/cron/notifications.py

