# cron/crontab
#crontab

#00	12  *   *   *	root 	*applications/medibox/cron/notifications.py
*   *   *   *   *	root	*applications/medibox/cron/test.py
